# hindsight.service
[Unit]
Description=MemoMind - Local Agent Memory Server
After=network.target

[Service]
Type=simple
User=hindsight
Environment=HF_ENDPOINT=https://hf-mirror.com
ExecStart=/opt/hindsight-env/bin/python3 /opt/hindsight-env/serve.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
