# memomind.service
[Unit]
Description=MemoMind - Local Agent Memory Server
After=network.target

[Service]
Type=simple
User=memomind
ExecStart=/opt/memomind-env/start-memomind.sh
Restart=on-failure
RestartSec=30
TimeoutStartSec=660

[Install]
WantedBy=multi-user.target
